# files/dehydrated.timer
[Unit]
Description=Run dehydrated regularly

[Timer]
Persistent=true
OnCalendar=daily
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target

